# examples/mcp/com.agentkitchen.chatgpt-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentkitchen.chatgpt-mcp</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>set -a; if [ -f "$HOME/.agent-kitchen/com.agentkitchen.chatgpt-mcp.env" ]; then . "$HOME/.agent-kitchen/com.agentkitchen.chatgpt-mcp.env"; fi; set +a; exec "__AGENT_KITCHEN_ROOT__/scripts/agentkitchen-mcp.sh" --http --host 0.0.0.0 --port 8765</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>AGENT_KITCHEN_ROOT</key>
        <string>__AGENT_KITCHEN_ROOT__</string>
        <key>KITCHEN_MCP_PUBLIC_BASE_URL</key>
        <string>https://kitchen.example</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/agentkitchen-chatgpt-mcp.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/agentkitchen-chatgpt-mcp.log</string>
</dict>
</plist>
